# back/mintzai-httpd/.service
[Unit]
Description=mintzai
After=network.target

[Service]
Environment="AHOLAB_CSV_FILE_PATH=/home/aholab/aitor/file.csv"
Environment="AHOLAB_AUDIO_FILES=/home/aholab/aitor/front/audios/"
Environment="AHOLAB_ACCEPTED_ORIGINS=https://aholab.ehu.eus/users/aitor/test.html"
Type=simple
Restart=always
RestartSec=30s
ExecStart=/home/aholab/aitor/back/httpd/httpd
StandardOutput=file:/var/log/mintzai.log
StandardError=file:/var/log/mintzai.err

[Install]
WantedBy=multi-user.target
